# texlive/texlive/files/etc/profile.d/texlive.sh
# Begin /etc/profile.d/texlive.sh

pathprepend /opt/texlive/bin           PATH

# End /etc/profile.d/texlive.sh
