# dnsmasq/dnsmasq/dnsmasq.etc.files
/etc/dnsmasq.conf
/etc/logrotate.d/dnsmasq
/etc/conf.d/dnsmasq
/etc/init.d/dnsmasq-dhcp
/etc/init.d/dnsmasq
/etc/dnsmasq.d/.delete_me
