# ntp/ntp/files/etc/init.d/ntpd
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

NTPD_OPTS="${NTPD_OPTS:--c /etc/ntp.conf}"

command=/usr/sbin/ntpd
pidfile=/run/ntpd.pid
command_args="$NTPD_OPTS"

depend() {
    use net dns logger
    after ntp-client
}

start_pre() {
    if [ ! -f /etc/ntp.conf ]; then
        eerror "Missing /etc/ntp.conf"
        eerror "See ntpd.conf(5)"
        return 1
    fi
    return 0
}
