# atom/atom/files/etc/profile.d/atom.sh
# Begin /etc/profile.d/atom.sh

# Adjust PATH
pathappend /opt/atom/bin

# End /etc/profile.d/atom.sh
