# openrc/openrc/files/etc/init.d/vbox-vmsvga
#!/sbin/openrc-run
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need xdm
}

start() {
	if [ -f /usr/bin/VBoxClient ]; then
		if [ -f /var/run/lightdm.pid ] || [ -f /var/run/xdm.pid ]; then
			/usr/bin/VBoxClient --vmsvga
		fi
	fi;
}

